# deploy/supervisor/medca-worker.conf
; MedcaHealthcare queue worker — deploy to /etc/supervisor/conf.d/medca-worker.conf
; Handles: ScoreLeadPayloadJob, SendOutboundWebhookJob, RefreshAiPulseSnapshotJob
;
; After deploy:
;   sudo supervisorctl reread
;   sudo supervisorctl update
;   sudo supervisorctl status medca-worker:*

[program:medca-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/medcahealthcare/artisan queue:work redis --sleep=3 --tries=3 --max-time=3600 --timeout=90
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/medcahealthcare/storage/logs/worker.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stopwaitsecs=3600
